# files/udev-cleanup.timer
[Unit]
Description=Clean up udev database

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
